# scripts/com.directstay.wake-doorbell.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.directstay.wake-doorbell</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/agents/.openclaw/workspace/directstay/scripts/run-directstay-wake-doorbell.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/agents/.openclaw/logs/directstay-wake-doorbell.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/agents/.openclaw/logs/directstay-wake-doorbell.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/agents/.openclaw/workspace/directstay</string>
</dict>
</plist>
