# scripts/com.directstay.platform-lead-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.directstay.platform-lead-worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/agents/.openclaw/workspace/villa-la-percha/scripts/run-directstay-platform-lead-worker.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/agents/.openclaw/logs/directstay-platform-lead-worker.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/agents/.openclaw/logs/directstay-platform-lead-worker.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/agents/.openclaw/workspace/villa-la-percha</string>
</dict>
</plist>
